# archive_prod_folders_20251027/production/redhat-deployment/prod/systemd/hrm-frontend.service
[Unit]
Description=Skyraksys HRM Frontend Service
Documentation=https://github.com/your-org/skyraksys-hrm
After=network-online.target hrm-backend.service
Wants=network-online.target hrm-backend.service

[Service]
Type=simple
User=hrmapp
Group=hrmapp
WorkingDirectory=/opt/skyraksys-hrm/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
# Note: CRA embeds REACT_APP_* at build time. Runtime env here won't change API URL.
# Keep build-time REACT_APP_API_URL set during `npm run build`.

# Load environment variables from file (if needed)
EnvironmentFile=-/opt/skyraksys-hrm/.env

# Start command for production (ESM-safe): pin serve to CommonJS version
# Prevent ERR_REQUIRE_ESM by invoking CLI and pinning serve@14 explicitly
ExecStart=/usr/bin/npx --yes serve@14 -s build -l 3000

# Alternative: Use npm start for development mode
# ExecStart=/usr/bin/npm start

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Health check
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=append:/var/log/skyraksys-hrm/frontend.log
StandardError=append:/var/log/skyraksys-hrm/frontend-error.log
SyslogIdentifier=hrm-frontend

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# Memory and CPU limits (frontend typically needs less resources)
MemoryMax=512M
CPUQuota=100%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# File system access
ReadWritePaths=/opt/skyraksys-hrm/frontend
ReadWritePaths=/var/log/skyraksys-hrm
ReadWritePaths=/tmp

# Network access
PrivateNetwork=false
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Capabilities (minimal required)
CapabilityBoundingSet=
AmbientCapabilities=

# System calls
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @privileged @reboot @swap @raw-io

[Install]
WantedBy=multi-user.target