# redhatprod/systemd/postgresql-skyraksys.service
[Unit]
Description=PostgreSQL database server for Skyraksys HRM
Documentation=man:postgres(1)
After=network-online.target
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=notify
User=postgres
Group=postgres

# PostgreSQL configuration for Skyraksys HRM
Environment=PGDATA=/var/lib/pgsql/data
Environment=PGUSER=postgres
Environment=PGDATABASE=skyraksys_hrm

# Service startup
ExecStart=/usr/bin/postgres -D ${PGDATA}
ExecReload=/bin/kill -HUP $MAINPID

# Service management
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/pgsql

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Restart policy
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target