# redhat/systemd/skyraksys-hrm.service
[Unit]
Description=SkyrakSys HRM Application
Documentation=https://github.com/Otyvino/skyrakskys_hrm
After=network.target postgresql-15.service
Wants=postgresql-15.service

[Service]
Type=simple
User=USER_PLACEHOLDER
Group=USER_PLACEHOLDER
WorkingDirectory=APP_DIR_PLACEHOLDER
Environment=NODE_ENV=production
Environment=PORT=8080
ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=60
RestartSec=5
Restart=on-failure
RestartPreventExitStatus=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=APP_DIR_PLACEHOLDER/logs APP_DIR_PLACEHOLDER/uploads /tmp

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=append:/var/log/skyraksys_hrm/application.log
StandardError=append:/var/log/skyraksys_hrm/error.log
SyslogIdentifier=skyraksys-hrm

[Install]
WantedBy=multi-user.target
